# dist/mastodon-streaming.service
[Unit]
Description=truecolors-streaming
After=network.target
Wants=truecolors-streaming@4000.service

[Service]
Type=oneshot
ExecStart=/bin/echo "truecolors-streaming exists only to collectively start and stop truecolors-streaming@ instances, shimming over the migration to templated truecolors-streaming systemd unit"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
